# backup/crontab
0 * * * * /bin/bash /backup.sh --tag=automated
